# supervisor.conf
[supervisord]
nodaemon=true

[program:rc.local]
command=/etc/rc.local
autorestart=false
priority=0
startsecs=0

[program:riak]
command=/usr/sbin/riak console
autorestart=false
